# roles/photos_archive/templates/mnt-photosarchive.automount
[Unit]
Description=Auto-mount photos archive share
After=network-online.target
Wants=network-online.target

[Automount]
Where=/mnt/photosarchive
TimeoutIdleSec=5min

[Install]
WantedBy=multi-user.target
